# systemd/irrigation.timer
[Unit]
Description=Triggers the irrigation service

[Timer]
OnCalendar=*-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target